# supervisord.conf
[supervisord]
nodaemon=true

[program:armviz]
command=cd /tmp/AzureResourceVisualizer/ && gulp serve
